# deploy/jevis.service
[Unit]
Description=JEVis Data Management Platform
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=jevis
Group=jevis
WorkingDirectory=/opt/jevis

ExecStart=/usr/bin/java \
  -Xms512m -Xmx2g \
  -jar target/JEVis4_Claude-1.0-SNAPSHOT.jar \
  --spring.config.additional-location=file:/etc/jevis/application.properties

Restart=always
RestartSec=10
SuccessExitStatus=143

StandardOutput=journal
StandardError=journal

# Sicherheits-Härtung
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=/var/log/jevis /opt/jevis/jte-classes
ReadOnlyPaths=/etc/jevis

[Install]
WantedBy=multi-user.target
